# .system-files/etc/systemd/system/docker.service.d/override.conf
[Service]
ExecStartPre=/usr/local/libexec/docker-wait-for-network
